# env/chef/cookbooks/upstart/templates/default/ripple_gateway-webapp.conf.erb
start on starting ripple_gateway
stop on stopping ripple_gateway
