# env/chef/cookbooks/upstart/templates/default/ripple_gateway.conf.erb
description "ripple_gateway"
author "Steven Zeiler <zeiler.steven@gmail.com>"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 99 5

script
  export HOME="/home/ubuntu"

  cd /var/www/ripple_gateway/current
  exec /usr/bin/node processes/webapp.js >> node.log 2>&1
end script
