# env/chef/cookbooks/upstart/templates/default/ripple_gateway-outgoing_ripple_payments-1.conf.erb
start on starting ripple_gateway-outgoing_ripple_payments
stop on stopping ripple_gateway-outgoing_ripple_payments
respawn

exec su - ubuntu -c 'cd /var/www/ripple_gateway/current; export PORT=5100; node processes/outgoing.js >> /var/log/ripple_gateway/outgoing_ripple_payments-1.log 2>&1'
